# scripts/chf.service.in
[Unit]
Description=3GPP 5G Core Charging Function (CHF)
Documentation="https://github.com/sigscale/chf"
After=epmd.service epmd.socket

[Service]
User=otp
Group=otp
WorkingDirectory=/home/otp
RuntimeDirectory=chf
RuntimeDirectoryMode=0750
EnvironmentFile=/etc/default/chf
LimitNOFILE=4096
ExecStart=@ERLANG_ROOT_DIR@/bin/run_erl -daemon /run/chf/ log \
                'ERL_LIBS=lib exec ${ROOTDIR}/bin/start_erl ${ROOTDIR} releases \
                ${START_ERL_DATA} -boot_var OTPHOME . +K true +A 32 +Bi -sname ${NODENAME}'
GuessMainPID=no
Type=forking
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
